# package/suderra-watchdog/suderra-watchdog.service
# SPDX-FileCopyrightText: 2026 Suderra OS contributors
# SPDX-License-Identifier: Apache-2.0
#
# Donanım watchdog'unun TEK sahibi (NEW-4). Kernel /dev/watchdog'a tek opener'a
# izin verir; bu daemon cihaz genelinde beslemeyi ve kademeli eskalasyonu yönetir.
# İzlenecek servis (ör. suderra-agent.service) komisyonlamada bir drop-in ile
# SUDERRA_WATCHDOG_HEALTH_UNIT üzerinden verilir — varsayılan BOŞtur ki agent
# kurulmamış bir cihaz sağlıklıyken restart/reboot eskalasyonuna girmesin.
# QEMU/dev'de donanım watchdog yoksa daemon soft moda düşer (REQUIRE_HW=1 ile
# donanım şartı zorlanabilir; G4 donanım lane'i bunu sıkılaştırır).

[Unit]
Description=Suderra hardware watchdog owner (feed + tiered escalation)
Documentation=https://github.com/Okan-wqm/Suderra-OS
# Donanım watchdog'unu MÜMKÜN OLDUĞUNCA ERKEN beslemeye başla: erken boot'ta
# takılan bir servis, geç başlayan bir watchdog tarafından yakalanamaz.
# (Sağlık-probe HEALTH_UNIT default BOŞ olduğundan agent'ı beklemesi gerekmez.)
DefaultDependencies=no
After=local-fs.target systemd-modules-load.service
Before=sysinit.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/suderra-watchdog
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX
LockPersonality=true
RemoveIPC=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

DevicePolicy=closed
DeviceAllow=/dev/watchdog rw

StandardInput=null
StandardOutput=journal
StandardError=journal
SyslogIdentifier=suderra-watchdog

[Install]
WantedBy=sysinit.target
