# board/suderra/common/rootfs-overlay/etc/systemd/system/suderra-firstboot.service
#
# Suderra OS — firstboot servisi
#
# İlk boot'ta çalışır:
#   1. /data partition and bind mounts are prepared by suderra-data.service
#   2. /var/lib/suderra runtime directories are initialized
#   3. /etc/machine-id üret (Buildroot empty bırakır)
#   4. Edge Agent için config skeleton oluştur (yoksa)
#   5. Tamamlandı flag'i koy: /var/lib/suderra/.firstboot-done
#
# Bu unit `suderra-firstboot` Rust binary'sini çağırır (Faz 2). Faz 1'de
# inline shell ile minimum işlevsellik sağlanır.
#
# ConditionPathExists ile bir kere çalışır (idempotent).
#

[Unit]
Description=Suderra OS firstboot initialization
Documentation=https://github.com/Okan-wqm/suderra-os
DefaultDependencies=no
After=local-fs.target systemd-tmpfiles-setup.service suderra-data.service
Before=suderra-agent.service
ConditionPathExists=!/var/lib/suderra/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes

# Faz 2'de Rust binary: /usr/bin/suderra-firstboot
# Faz 1 minimum: inline shell (data partition mkfs + dizin hazırlama)
ExecStart=/bin/sh -c ' \
    set -e; \
    echo "==> Suderra firstboot başlatılıyor"; \
    \
    # 1. /etc/machine-id (Buildroot empty bırakır) \
    if [ ! -s /etc/machine-id ]; then \
        systemd-machine-id-setup; \
    fi; \
    \
    # 2. Mutable state/config and provisioning directories \
    mkdir -p /var/lib/suderra/config /var/lib/suderra/edge /var/log/suderra /etc/suderra /var/spool/suderra-provision/inbox; \
    chown root:root /var/lib/suderra /var/spool/suderra-provision 2>/dev/null || true; \
    chown 200:200 /var/lib/suderra/config /var/lib/suderra/edge /var/log/suderra 2>/dev/null || true; \
    chown 201:201 /var/spool/suderra-provision/inbox 2>/dev/null || true; \
    chown root:root /etc/suderra 2>/dev/null || true; \
    chmod 0750 /var/lib/suderra /var/lib/suderra/config /var/lib/suderra/edge /var/log/suderra 2>/dev/null || true; \
    chmod 0750 /var/spool/suderra-provision 2>/dev/null || true; \
    chmod 0700 /var/spool/suderra-provision/inbox 2>/dev/null || true; \
    chmod 0755 /etc/suderra 2>/dev/null || true; \
    \
    # 4. Config skeleton (mutable, persisted under /var/lib) \
    if [ ! -f /var/lib/suderra/config/config.yaml ]; then \
        printf "%s\n" \
          "device_id: \"00000000-0000-0000-0000-000000000000\"" \
          "device_code: \"UNACTIVATED\"" \
          "api_url: \"https://provisioning.example.com\"" \
          "mqtt:" \
          "  broker: null" \
          "  port: 8883" \
          "  username: null" \
          "  password: null" \
          "  tls:" \
          "    enabled: true" \
          "    verify_hostname: true" \
          "  topics: {}" > /var/lib/suderra/config/config.yaml; \
        chown 200:200 /var/lib/suderra/config/config.yaml 2>/dev/null || true; \
        chmod 0600 /var/lib/suderra/config/config.yaml; \
    fi; \
    \
    # 5. One-time provisioning password for the forced-command SSH user \
    if [ ! -f /var/spool/suderra-provision/credentials.env ]; then \
        pass="$(tr -dc A-HJ-NP-Za-km-z2-9 </dev/urandom | head -c 24)"; \
        hash="$(openssl passwd -6 "${pass}")"; \
        sed -i "s#^provision:[^:]*:#provision:${hash}:#" /etc/shadow; \
        printf "SUDERRA_PROVISION_USER=provision\nSUDERRA_PROVISION_PASSWORD=%s\n" "${pass}" > /var/spool/suderra-provision/credentials.env; \
        chown root:root /var/spool/suderra-provision/credentials.env 2>/dev/null || true; \
        chmod 0600 /var/spool/suderra-provision/credentials.env; \
        printf "Suderra OS provisioning user: provision\nPassword: %s\n" "${pass}" > /dev/console 2>/dev/null || true; \
    fi; \
    \
    # 6. Done flag \
    touch /var/lib/suderra/.firstboot-done; \
    echo "==> Suderra firstboot tamamlandı"; \
'

# Hardening (oneshot bile olsa)
NoNewPrivileges=true
ProtectSystem=strict
# /var/spool: provisioning inbox + tek-kullanımlık credentials.env buraya
# yazılır (yukarıda adım 2) — listeden düşerse set -e altında firstboot
# iptal olur ve cihaz hiç provision edilemez.
ReadWritePaths=/etc /var/lib /var/log /var/spool
ProtectHome=true
PrivateTmp=true
PrivateNetwork=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=suderra-firstboot

[Install]
WantedBy=sysinit.target
