# board/suderra/common/rootfs-overlay/etc/systemd/system/suderra-provision-worker.path
[Unit]
Description=Watch Suderra OS provisioning requests
Documentation=https://github.com/Okan-wqm/suderra-os

[Path]
PathModified=/var/lib/suderra/provision/request
MakeDirectory=true

[Install]
WantedBy=multi-user.target
