# board/suderra/common/rootfs-overlay/etc/systemd/system/suderra-data.service
[Unit]
Description=Suderra mutable data partition
Documentation=https://github.com/Okan-wqm/suderra-os
DefaultDependencies=no
Wants=local-fs.target
After=local-fs.target systemd-tmpfiles-setup.service
Before=suderra-firstboot.service suderra-agent.service nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/suderra-data-unlock

NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/data /var/lib/suderra /var/log/suderra /opt/suderra/edge/releases /dev/mapper

[Install]
WantedBy=sysinit.target
