# board/suderra/common/rootfs-overlay/etc/systemd/system/suderra-edge-agent.service
#
# Suderra Edge Agent — sertleştirilmiş systemd unit (Faz 2 placeholder)
#
# Bu unit Suderra Edge Agent v1.6.0 ile uyumlu olacak.
# Asıl unit dosyası package/suderra-edge-agent/ içinde Faz 2'de yazılır.
#

[Unit]
Description=Suderra Edge Agent (Aquaculture Industrial Edge)
Documentation=https://github.com/Okan-wqm/suderra-os
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/var/lib/suderra

[Service]
Type=notify
NotifyAccess=main

ExecStart=/usr/bin/suderra-edge-agent --config /etc/suderra/config.yaml

User=suderra-edge
Group=suderra-edge

# === Watchdog ===
WatchdogSec=60s
TimeoutStartSec=30s
TimeoutStopSec=90s
Restart=on-failure
RestartSec=5s

# === Resource limits ===
MemoryMax=256M
CPUQuota=75%
TasksMax=64
LimitNOFILE=1024
LimitCORE=0
LimitNPROC=128

# === Sertleştirme (systemd hardening) ===
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
DevicePolicy=closed
DeviceAllow=/dev/null rw
DeviceAllow=/dev/zero rw
DeviceAllow=/dev/urandom r
ReadWritePaths=/var/lib/suderra
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true

# === Capabilities ===
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# === Seccomp ===
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @debug @mount @swap @reboot @raw-io @keyring
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# === Network ===
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
IPAddressDeny=any
# IPAddressAllow listesi config'e bağlı — Faz 3'te detaylanır

# === Journal ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=suderra-edge-agent

[Install]
WantedBy=multi-user.target
