# package/suderra-rauc-config/suderra-rauc-mark-good.service
[Unit]
Description=Mark current Suderra OS RAUC slot good
Documentation=https://rauc.readthedocs.io/
After=multi-user.target boot.mount suderra-data.service rauc.service
Requires=boot.mount suderra-data.service
RequiresMountsFor=/data /boot
ConditionPathExists=/etc/rauc/system.conf
ConditionPathExists=/usr/sbin/suderra-rauc-health-gate
ConditionKernelCommandLine=|rauc.slot=A
ConditionKernelCommandLine=|rauc.slot=B

[Service]
Type=oneshot
ExecStart=/usr/sbin/suderra-rauc-mark-good
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/data /boot

[Install]
WantedBy=multi-user.target
